# configs/superd/services/autocutsel-primary.service
[Unit]
Description=automated xcutsel

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/autocutsel -selection PRIMARY
