# configs/superd/services/dunst.service
[Unit]
Description=lightweight notification daemon for x11 compositors

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/dunst
