# configs/superd/services/sxmo_hook_lisgd.service
[Unit]
Description=run user defined lisgd start hook

[Service]
Type=simple
Restart=always
ExecStart=sxmo_hook_lisgdstart.sh
