# configs/openrc/sxmo-setpermissions
#!/sbin/openrc-run

description="Set user permissions to aspects of the hardware kernel interface"

command="/usr/bin/sxmo_setpermissions.sh"
