# configs/external-services/mako.service
[Unit]
Description=lightweight notification daemon for Wayland compositors

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/mako
