# configs/external-services/vvmd.service
[Unit]
Description=lower level daemon that retrieves Visual Voicemail

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/vvmd
