# configs/services/geoclue-agent.service
[Unit]
Description=Proxy location data to map applications

[Service]
Type=simple
Restart=always
ExecStart=/usr/libexec/geoclue-2.0/demos/agent
