# configs/services/sxmo-x11-status.service
[Unit]
Description=periodically set the status bar content

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/sxmo_status_xsetroot.sh
