# configs/services/sxmo_autosuspend.service
[Unit]
Description=daemon for autosuspending

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/sxmo_autosuspend.sh
