# configs/services/sxmo_desktop_widget.service
[Unit]
Description=Sxmo: start desktop widget (ie. clock)

[Service]
Type=simple
Restart=always
ExecStart=sxmo_hook_desktop_widget.sh
