# configs/services/sxmo_battery_monitor.service
[Unit]
Description=Sxmo: monitors the battery

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/sxmo_battery_monitor.sh
