# configs/services/sxmo_swaybg.service
[Unit]
Description=The swaybg as a service

[Service]
Type=simple
Restart=always
ExecStart=swaybg -i "$SXMO_BG_IMG" -m fill
