# configs/services/sxmo_fehbg.service
[Unit]
Description=The feh bg as a service

[Service]
Type=simple
Restart=always
ExecStart=feh --bg-fill "$SXMO_BG_IMG"
