# configs/services/sxmo_notificationmonitor.service
[Unit]
Description=display notifications popup messages

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/sxmo_notificationmonitor.sh
