# configs/services/sxmo_conky.service
[Unit]
Description=Conky

[Service]
Type=simple
Restart=always
ExecStart=conky -c "$XDG_CONFIG_HOME"/sxmo/conky.conf
